# khal/root/etc/cron.d/vdirsyncer
*/15 * * * * dodo /bin/vdir_sync.sh >> /data/vdir/vdirsyncer.log 2>&1

