# khal/root/etc/cron.d/vdirsyncer
0 * * * * dodo /opt/python/bin/vdirsyncer sync >> /data/vdir/vdirsyncer.log 2>&1

